# root/etc/uci-defaults/luci-advancedplus
#!/bin/sh
chmod +x /etc/init.d/advancedplus

rm -f /tmp/luci-indexcache
exit 0
