# root/etc/uci-defaults/luci-advancedplus
#!/bin/sh
chmod +x /etc/init.d/advancedplus /usr/bin/advancedplusipk /usr/bin/uhttpd-watchdog


rm -f /tmp/luci-indexcache
exit 0
